# docs/systemd/conduit.service
[Unit]
Description=Conduit Protocol Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/conduit
Environment=NODE_ENV=production
Environment=CONDUIT_BIND=0.0.0.0
Environment=CONDUIT_HTTP_PORT=9087
Environment=CONDUIT_WS_PORT=9088
Environment=CONDUIT_RULES=config/rules.yaml
Environment=CONDUIT_CODECS_HTTP=true
Environment=CONDUIT_CODECS_WS=true
Environment=CONDUIT_TOKENS=dev-local
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
